# .config/autostart/02_stalonetray.desktop
[Desktop Entry]
Type=Application
Terminal=false
Exec=stalonetray
Name=Stalonetray
