# .config/autostart/02_screen.desktop
[Desktop Entry]
Exec=$HOME/bin/scripts/screens/home.sh
Name=Setup Displays
Type=Application
