# .config/autostart/04_nmapplet.desktop
[Desktop Entry]
Type=Application
Terminal=false
Exec=nm-applet
Name=Network Manager Applet
