# .config/autostart/04_sxhkd.desktop
[Desktop Entry]
Exec=/usr/bin/sxhkd -c /home/khazakar/.config/sxhkd/config
Name=sxhkd
Type=Application
