# .config/autostart/04_sxhkd.desktop
[Desktop Entry]
Exec=/bin/sh -c "/usr/bin/sxhkd -c $HOME/.config/sxhkd/config"
Name=sxhkd
Type=Application
