# .config/autostart/05_blueman.desktop
[Desktop Entry]
Type=Application
Terminal=false
Exec=blueman-applet
Name=Bluetooth Applet
