# .config/autostart/05_tint2.desktop
[Desktop Entry]
Exec=$HOME/bin/scripts/apps/tint2.sh
Name=tint2
Type=Application
