# .config/autostart/05_tint2.desktop
[Desktop Entry]
Exec=/bin/sh /home/khazakar/bin/scripts/apps/tint2.sh
Name=tint2
Type=Application
