# .config/autostart/06_blueman.desktop
[Desktop Entry]
Name=Blueman Applet
Icon=blueman
Exec=blueman-applet
Terminal=false
Type=Application
Categories=
NotShowIn=KDE;GNOME;
