# .config/autostart/06_pasystray.desktop
[Desktop Entry]
Type=Application
Terminal=false
Exec=pasystray
Name=PulseAudio System Tray
