# .config/autostart/07_cmst.desktop
[Desktop Entry]
Exec=/usr/bin/cmst
Name=QT5 Applet for Connman
Type=Application
