# .config/autostart/07_nmapplet.desktop
[Desktop Entry]
Exec=/usr/bin/nm-applet
Name=NetworkManager Applet
Type=Application
