# .config/autostart/07_xtile.desktop
[Desktop Entry]
Type=Application
Terminal=false
Exec=x-tile
Name=X-Tile
