# .config/autostart/08_kdeconnect.desktop
[Desktop Entry]
Exec=$HOME/bin/scripts/apps/kdeconnect.sh
Name=kdeconnect-indicator
Type=Application
