# .config/autostart/08_kdeconnect.desktop
[Desktop Entry]
Exec=/bin/sh /home/khazakar/bin/scripts/apps/kdeconnect.sh
Name=kdeconnect-indicator
Type=Application
