# .config/autostart/08_screens.desktop
[Desktop Entry]
Type=Application
Terminal=false
Exec=/home/khazakar/.screenlayout/home.sh
Name=Screen Layout at home
