# .config/autostart/09_background.desktop
[Desktop Entry]
Type=Application
Terminal=false
Exec=feh --bg-scale /home/khazakar/Obrazy/wall.png
Name=Set background
