# .config/autostart/11_picom.desktop
[Desktop Entry]
Exec=/usr/bin/picom -CGb
Name=picom
Type=Application
