# .config/autostart/12_locker.desktop
[Desktop Entry]
Exec=/bin/sh -c "cd $HOME/bin/scripts/apps ; ./locker.sh"
Name=locker.sh
Type=Application
