# .config/autostart/14_locker.desktop
[Desktop Entry]
Exec=$HOME/bin/scripts/apps/locker.sh
Name=locker.sh
Type=Application
